# bin/deploy/crontab
# vim: filetype=crontab
MAILTO=support-klarschiff@bfpi.de

# m h   dom mon dow  user      command
32  *    *   *   *  rails     $HOME/backoffice/bin/cron rails runner DeleteUnconfirmedIssuesJob.perform_now
37  *    *   *   *  rails     $HOME/backoffice/bin/cron rails runner DeleteUnconfirmedSupportersJob.perform_now
42  *    *   *   *  rails     $HOME/backoffice/bin/cron rails runner DeleteUnconfirmedAbuseReportsJob.perform_now
45  *    *   *   *  rails     $HOME/backoffice/bin/cron rails runner DeleteUnconfirmedCompletionsJob.perform_now
47  *    *   *   *  rails     $HOME/backoffice/bin/cron rails runner DeleteUnconfirmedPhotosJob.perform_now
52  *    *   *   *  rails     $HOME/backoffice/bin/cron rails runner DeleteAuthorsAfterDeadlineJob.perform_now
40  00   *   *   *  rails     $HOME/backoffice/bin/cron rails runner ArchiveClosedIssuesJob.perform_now
40  01   *   *   *  rails     $HOME/backoffice/bin/cron rails runner InformEditorialStaffOnIssuesJob.perform_now
00  02   *   *   *  rails     $HOME/backoffice/bin/cron rails runner CalculateAverageTurnaroundTimeJob.perform_now
00  07   *   *   *  rails     $HOME/backoffice/bin/cron rails runner InformOnDelegatedIssuesJob.perform_now
05  07   *   *   *  rails     $HOME/backoffice/bin/cron rails runner NotifyOnGroupAssignedResponsibilityJob.perform_now
05  10   *   *   *  rails     $HOME/backoffice/bin/cron rails runner NotifyOnIssuesInProcessJob.perform_now
10  10   *   *   *  rails     $HOME/backoffice/bin/cron rails runner NotifyOnClosedIssuesJob.perform_now
