# bin/deploy/crontab

0  02   *   *   *   rails     rails runner CalculateAverageTurnaroundTimeJob.perform_now
0  32   *   *   *   rails     rails runner UnconfirmedIssuesDeletionJob.perform_now
0  37   *   *   *   rails     rails runner UnconfirmedSupportsDeletionJob.perform_now
0  42   *   *   *   rails     rails runner UnconfirmedAbusesDeletionJob.perform_now
0  47   *   *   *   rails     rails runner UnconfirmedFotosDeletionJob.perform_now
0  40  00   *   *   rails     rails runner ArchiveClosedIssuesJob.perform_now
0  40  01   *   *   rails     rails runner InformEditorialStaffOnIssuesJob.perform_now
0  00  07   *   *   rails     rails runner InformOnDelegatedIssuesJob.perform_now
0  05  10   *   *   rails     rails runner NotificationOnInProcessJob.perform_now
0  10  10   *   *   rails     rails runner NotificationClosedJob.perform_now
