# bin/deploy/crontab

0  02   *   *   *   rails     rails runner CalculateAverageTurnaroundTimeJob.perform_now
